# systemd/homelab-host-backup.timer
[Unit]
Description=Daily Proxmox host configuration backup

[Timer]
OnCalendar=*-*-* 01:15:00
Persistent=true
RandomizedDelaySec=10m

[Install]
WantedBy=timers.target
