# systemd/homelab-host-backup.service
[Unit]
Description=Back up Proxmox host configuration to USB backup drive
Wants=mnt-proxmox\x2dusb\x2dbackup.mount
After=mnt-proxmox\x2dusb\x2dbackup.mount

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/homelab-host-backup
